# ansible/rdm.service
[Unit]
Description=rdm

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu
ExecStart=/usr/local/bin/rdm

[Install]
WantedBy=multi-user.target